# misc/cobalt.in
#!/bin/sh
#
# Cobalt
#
# chkconfig: 2345 18 09
# description: Cobalt resource manager/scheduler
#

prefix=@prefix@
exec_prefix=@exec_prefix@
CONFIGPATH=@sysconfdir@
CONFIGFILE="${CONFIGPATH}/cobalt.conf"
DAEMONDIR=@sbindir@
PIDPATH=@PATH_TO_PIDFILE@

# Configure DB2 environment
. /discovery/db.src

case "$1" in
  start)
    echo -n "Starting Process Management: "
    "${DAEMONDIR}"/bgpm.py -C "${CONFIGFILE}" --daemon "${PIDPATH}"/bgpm.pid
    echo "bgpm.py"
    echo -n "Starting Queue Manager: "
    "${DAEMONDIR}"/cqm.py -C "${CONFIGFILE}" --daemon "${PIDPATH}"/cqm.pid
    echo "cqm.py"
    echo -n "Starting Scheduler: "
    "${DAEMONDIR}"/bgsched.py -C "${CONFIGFILE}" --daemon "${PIDPATH}"/bgsched.pid
    echo "bgsched.py"
    touch /var/lock/subsys/cobalt
    ;;
  stop)
    if [ -f "${PIDPATH}/bgsched.pid" ] ; then
	echo -n "Stopping Scheduler: "
	kill -INT `cat "${PIDPATH}"/bgsched.pid`
	echo done
    else
	echo "Scheduler not running"
    fi
    if [ -f "${PIDPATH}/cqm.pid" ] ; then
	echo -n "Stopping Queue Manager: "
	kill -INT `cat "${PIDPATH}"/cqm.pid`
	echo done
    else
	echo "Cobalt Queue Manager not running"
    fi
    if [ -f "${PIDPATH}/bgpm.pid" ] ; then
	echo -n "Stopping Process Management: "
	kill -INT `cat "${PIDPATH}"/bgpm.pid`
	echo done
    else
	echo "Cobalt Blue Gene Process Manager not running"
    fi
    rm -f "${PIDPATH}"/bgsched.pid "${PIDPATH}"/cqm.pid "${PIDPATH}"/bgpm.pid
    rm -f /var/lock/subsys/cobalt
    ;;
  restart)
    $0 stop
    $0 start
    ;;
  force-reload)
    $0 stop
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|force-reload}"
    exit 1
esac

exit 0



